# ePaperCalendar/Calendar.service
[Unit]
Description=ePaperCalendar
After=network.target

[Service]
ExecStart=/path/to/your/script.sh
User=root

[Install]
WantedBy=multi-user.target
